# pi/systemd/vritti-kiosk.service
[Unit]
Description=Vritti Kiosk (fullscreen mandala)
After=ai-runtime.service graphical.target
Wants=ai-runtime.service
PartOf=graphical.target

[Service]
Type=simple
User=pi
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/pi/.Xauthority
ExecStartPre=/bin/sleep 5
ExecStart=/usr/bin/chromium-browser \
  --kiosk \
  --noerrdialogs \
  --disable-infobars \
  --disable-session-crashed-bubble \
  --disable-restore-session-state \
  --disable-features=TranslateUI \
  --check-for-update-interval=31536000 \
  --no-first-run \
  --start-fullscreen \
  --use-fake-ui-for-media-stream \
  --enable-media-stream \
  --autoplay-policy=no-user-gesture-required \
  --use-gl=egl \
  http://127.0.0.1:8000/
Restart=on-failure
RestartSec=10

[Install]
WantedBy=graphical.target
